# container/inittab
::respawn:/sbin/syslogd -n -O -
::respawn:/opt/start
